# pxa/base-files/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K shutdown
ttyS2::askfirst:/bin/ash --login
tty1::respawn:/bin/ash --login
